# build/docker-release/languagetool.supervisord.conf
[program:languagetool]
command=java -cp languagetool-server.jar org.languagetool.server.HTTPServer --port 8081 --allow-origin "*" --public
directory=/srv/inclusify/api
autostart=true
autorestart=true
startretries=5
numprocs=1
username=inclusify
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/%(program_name)s_stdout.log
stdout_logfile_maxbytes=10MB
